# bash/profile.d/999-node
# Load nodenv unless it loaded.
# append_path_unless_exist "$HOME/.nodenv/bin"
#
# if command -v nodenv > /dev/null; then
#   eval "$(nodenv init -)"
# fi
