# bash/profile.d/999-fzf
if command -v fzf > /dev/null; then
  export FZF_DEFAULT_COMMAND='git ls-files'
fi

